# ooproxy.conf
description "OpenOffice Proxy"
author "martin.reisenhofer@funkring.net"

start on runlevel [2345]
stop on runlevel [!2345]

exec /opt/ooproxy/ooproxy.py

respawn
respawn limit 10 90